# build/proxy/supervisor.conf
[supervisord]
nodaemon = true

[program:nginx]
command = usr/sbin/nginx -g "daemon off;"
user = root
autostart = true
autorestart = true
priority = 10

[program:cron]
command = /usr/sbin/cron -f -L 15
user = root
autostart = true
autorestart = true
priority = 10