# build/app/config/supervisor.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
user=root
autostart=true

[program:memcached]
command=/usr/bin/memcached -v -u root
user=root
autostart=true

[program:php-fpm]
command=/usr/sbin/php5-fpm -F
user=root
autostart=true