# build/app/supervisor.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron
user = root
autostart = true

[program:memcached]
command=/usr/bin/memcached -Lv
user=www-data
autostart=true

[program:php-fpm]
command=/usr/sbin/php-fpm7.0 -F
user=www-data
autostart=true